# iot-gui-webhook.service
[Unit]
Description=GitHub Webhook Listener for AWS IoT Pub/Sub GUI
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Pideployment
Environment="PATH=/home/pi/Pideployment/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="DISPLAY=:0"
Environment="WEBHOOK_PORT=9000"
Environment="WEBHOOK_HOST=0.0.0.0"
Environment="GIT_BRANCH=main"
# Webhook secret will be loaded from .webhook_secret file
# The setup script will replace WEBHOOK_SECRET_PLACEHOLDER with actual secret

# Use bash to load secret from file and start webhook listener
ExecStart=/bin/bash -c 'source /home/pi/Pideployment/venv/bin/activate && export WEBHOOK_SECRET=$(cat /home/pi/Pideployment/.webhook_secret 2>/dev/null || echo "change-me-to-a-strong-secret-key") && exec /home/pi/Pideployment/venv/bin/python /home/pi/Pideployment/webhook_listener.py'

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iot-gui-webhook

[Install]
WantedBy=multi-user.target
